# src/crontab/pakfire-build-service
# Send queued emails once a minute
* * * * *	pakfire	pakfire-build-service process-message-queue &>/dev/null

# Check build dependencies
* * * * *	pakfire	pakfire-build-service check-build-dependencies &>/dev/null

# Remaster repositories
*/5 * * * *	pakfire	pakfire-build-service remaster-repositories &>/dev/null

# Cleanup repositories
5 * * * *	pakfire	pakfire-build-service cleanup-repositories &>/dev/null

# Pull sources
*/5 * * * *	pakfire	pakfire-build-service pull-sources &>/dev/null

# Dist
*/5 * * * *	pakfire	pakfire-build-service dist &>/dev/null

# Send updates to Bugzilla
*/5 * * * *	pakfire	pakfire-build-service send-bug-updates &>/dev/null

# Create test jobs
*/15 * * * *	pakfire	pakfire-build-service create-test-jobs &>/dev/null

# Cleanup files
*/5 * * * *	pakfire	pakfire-build-service cleanup-files &>/dev/null

# Cleanup timed-out uploads
0 */6 * * *	pakfire	pakfire-build-service cleanup-uploads &>/dev/null

# Cleanup expired sessions
0 0 * * *	pakfire	pakfire-build-service cleanup-sessions &>/dev/null

# Run mirror check
*/30 * * * *	pakfire	pakfire-build-service check-mirrors &>/dev/null
