# src/crontab/pakfire-build-service
# Send queued emails once a minute
* * * * *	nobody	pakfire-build-service process-message-queue &>/dev/null
