# src/crontab/pakfire-build-service
# Send queued emails once a minute
* * * * *	nobody	pakfire-build-service process-message-queue &>/dev/null

# Send updates to Bugzilla
*/5 * * * *	nobody	pakfire-build-service send-bug-updates &>/dev/null

# Cleanup timed-out uploads
0 */6 * * *	nobody	pakfire-build-service cleanup-uploads &>/dev/null

# Cleanup expired sessions
0 0 * * *	nobody	pakfire-build-service cleanup-sessions &>/dev/null

# Run mirror check
*/30 * * * *	nobody	pakfire-build-service check-mirrors &>/dev/null
