# src/systemd/pakfire-web-1.service.in
[Unit]
Description=Pakfire Web UI #1
After=network.target

[Service]
ExecStart=@bindir@/pakfire-web --port=9001
User=nobody
Group=nogroup

[Install]
WantedBy=multi-user.target
