# src/systemd/pakfire-web-2.service.in
[Unit]
Description=Pakfire Web UI #2
After=network.target

[Service]
ExecStart=@bindir@/pakfire-web --port=9002
User=nobody
Group=nogroup

[Install]
WantedBy=multi-user.target
