# src/systemd/pakfire-hub-2.service.in
[Unit]
Description=Pakfire Hub #1
After=network.target

[Service]
ExecStart=@bindir@/pakfire-hub --port=8002
User=nobody
Group=nogroup

[Install]
WantedBy=multi-user.target
